# deploy/ubuntu/twag-telegram-agent-boston@.service
[Unit]
Description=TWAG Telegram agent, Boston (%i)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=%i
Group=%i
WorkingDirectory=__APP_DIR__
EnvironmentFile=/etc/twag/twag.env
Environment=TWAG_CITY=boston
Environment=TELEGRAM_AGENT_LOCK_FILE=/run/twag-boston/telegram-agent.lock
RuntimeDirectory=twag-boston
RuntimeDirectoryMode=0750
ExecStart=__APP_DIR__/.venv/bin/twag-telegram-agent
Restart=always
RestartSec=5
KillSignal=SIGINT
TimeoutStopSec=30

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/run/twag-boston /var/log/twag

[Install]
WantedBy=multi-user.target
