# deploy/ubuntu/twag-telegram-agent@.service
[Unit]
Description=TWAG Telegram agent (%i)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=%i
Group=%i
WorkingDirectory=__APP_DIR__
EnvironmentFile=/etc/twag/twag.env
RuntimeDirectory=twag
RuntimeDirectoryMode=0750
ExecStart=__APP_DIR__/.venv/bin/twag-telegram-agent
Restart=always
RestartSec=5
KillSignal=SIGINT
TimeoutStopSec=30

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/run/twag /var/log/twag

[Install]
WantedBy=multi-user.target
